# debian/postinst
if [ -f /etc/init.d/gluu-gateway ]; then
    /etc/init.d/konga stop > /dev/null 2>&1
    /etc/init.d/kong stop > /dev/null 2>&1
    /etc/init.d/oxd-server stop > /dev/null 2>&1
    /usr/sbin/update-rc.d kong remove > /dev/null 2>&1
    /usr/sbin/update-rc.d konga remove > /dev/null 2>&1
    /usr/sbin/update-rc.d oxd-server remove > /dev/null 2>&1
    /usr/sbin/update-rc.d gluu-gateway defaults > /dev/null 2>&1
fi
if [ -f /lib/systemd/system/gluu-gateway.service ]; then
     systemctl stop konga > /dev/null 2>&1
     systemctl stop kong > /dev/null 2>&1
     systemctl stop oxd-server > /dev/null 2>&1
fi
if [ -f /opt/gluu-gateway/setup/setup-gluu-gateway.py ]; then
     chmod +x /opt/gluu-gateway/setup/setup-gluu-gateway.py > /dev/null 2>&1
fi
