# debian/postinst
/usr/sbin/update-rc.d kong defaults
/usr/sbin/update-rc.d konga defaults
/etc/init.d/gluu-gateway stop 2>/dev/null
