# debian/prerm
#!/bin/sh
set -e
if [ -x "/etc/init.d/gluu-gateway" ] || [ -e "/etc/init/gluu-gateway.conf" ]; then
	invoke-rc.d gluu-gateway stop  2> /dev/null || exit $?
fi
